# ops/systemd/frpclient-postgres-verify.service
[Unit]
Description=FRP Client PostgreSQL backup restore verification
Wants=docker.service network-online.target
After=docker.service network-online.target frpclient-postgres-backup.service

[Service]
Type=oneshot
User=root
WorkingDirectory=/var/www/FRPclient
Environment=PROJECT_DIR=/var/www/FRPclient
Environment=BACKUP_DIR=/var/backups/frpclient/postgres
Environment=VERIFY_DB_PREFIX=frpclient_restore_verify
Environment=MIN_TABLE_COUNT=10
ExecStart=/bin/sh /var/www/FRPclient/ops/backup/postgres_verify_backup.sh
